# image/etc/systemd/system/shed-aws-proxy.service
[Unit]
Description=shed AWS Credential Proxy
After=network.target shed-agent.service
Wants=shed-agent.service

[Service]
Type=simple
User=shed
Group=shed
RuntimeDirectory=shed-extensions
AmbientCapabilities=CAP_NET_BIND_SERVICE
ExecStart=/usr/local/bin/shed-aws-proxy \
  --port 499 \
  --publish-url http://127.0.0.1:498/v1/publish
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
