# image/etc/systemd/system/shed-ssh-agent.service
[Unit]
Description=shed SSH Agent (credential proxy)
After=network.target shed-agent.service
Wants=shed-agent.service

[Service]
Type=simple
User=shed
Group=shed
RuntimeDirectory=shed-extensions
ExecStart=/usr/local/bin/shed-ssh-agent \
  --sock /run/shed-extensions/ssh-agent.sock \
  --publish-url http://127.0.0.1:498/v1/publish
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
